# roles/collectd/files/collectd.conf
[program:collectd]
command=/usr/sbin/collectd -f
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
priority=100
