# roles/jboss/files/jboss.conf
[program:jboss]
command=/indicee/bin/start_jboss.sh
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
priority=100
